# src/salt/base/state/network/wwan/qmi/qmi-manager.service
[Unit]
Description=QMI Network Connection Manager Service
After=network.target

[Service]
ExecStart=/usr/bin/qmi-manager run
ExecStopPost=/usr/bin/qmi-manager down
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target