# src/salt/base/state/network/wwan/qmi/udev.rules
ACTION=="add", SUBSYSTEM=="usbmisc", DRIVER=="", KERNEL=="cdc-wdm0", RUN+="/bin/sh -c 'echo Y > /sys/class/net/wwan0/qmi/raw_ip'"